# hhvm/supervisor-config/supervisor-queue.conf
[program:queue]
command = /scripts/supervisor-queue.sh
numprocs = 1
autostart = true
autorestart = true