# hhvm/supervisor-config/supervisor-hhvm.conf
[program:hhvm]
command = hhvm --mode server
numprocs = 1
autostart = true
autorestart = true